# systemd/subserver.service
[Unit]
Description=VPN Subscription Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/subserver/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
