# systemd/vpnbot.service
[Unit]
Description=VPN Telegram Bot
After=network.target xray.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpnbot
ExecStart=/usr/bin/python3 /opt/vpnbot/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
